# deploy/ai-employee-cloud.service
[Unit]
Description=AI Employee Cloud Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ai-employee
Group=ai-employee
WorkingDirectory=/opt/ai-employee
ExecStart=/opt/ai-employee/.venv/bin/python scripts/scheduler.py --mode cloud --vault /opt/ai-employee/vault
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
EnvironmentFile=/etc/ai-employee/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/ai-employee/vault /opt/ai-employee/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
